# data/session-settings-service.conf.in
description "Session settings service"
author "William Hua <william.hua@canonical.com>"

start on started dbus
stop on stopping dbus

respawn

exec @SERVICE_INSTALL_DIR@/session-settings-service
